# x11-test/x11.conf
[program:vnc]
command=vncserver "%(ENV_DISPLAY)s" -geometry "%(ENV_DISPLAY_SCREEN)s" -depth 16
autorestart=true
environment=USER="root",XMODIFIERS="@im=fcitx"
