# x11-client/files/etc/supervisor/conf.d/x11-client.conf
[program:dbus]
command=service dbus start

[program:fcitx]
command=fcitx