# chrome/files/etc/supervisor/conf.d/chrome.conf
[program:chrome]
# --disable-gpu
user=
command=google-chrome --user-data-dir=/data/.chrome --homedir=/data
autorestart=true