# chrome/chrome.conf
[program:chrome]
# fixing gpu
command=chromium-browser --disable-gpu --user-data-dir=/data/user-data --ppapi-flash-path=/usr/lib/PepperFlash
autorestart=true